# systemd/wifi-reconnect.service
[Unit]
Description=WiFi Check Service

[Service]
Type=simple
ExecStart=$HOME/wifi-reconnect.sh
