# systemd/wifi-reconnect.timer
[Unit]
Description=WiFi Check Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=2min
Unit=wifi-reconnect.service

[Install]
WantedBy=timers.target
