# systemd/wifi-reconnect.service
[Unit]
Description=WiFi Check Service

[Service]
Type=simple
ExecStart=/home/david/wifi-reconnect.sh
